# conf/nas/dropbox/dropboxd.conf
# dropbox

description "Dropbox daemon"
author "Miphreal Adler"

start on filesystem and stopped networking
stop on runlevel [016]

respawn
respawn limit 5 30

env uid=nas-user
env gid=nas-user
env umask=007

exec start-stop-daemon -S -c $uid:$gid -k $umask -x /usr/bin/dropboxd >> /var/log/dropbox/daemon.log